# templates/supervisor/quakelive.conf.j2
[program:quakelive]
command=/home/{{ quakelive_user }}/quakelive/run_server_x64_minqlx.sh
user={{ quakelive_user }}
numprocs=1
autorestart=true
